# assets/supervisord.conf
[supervisord]
nodaemon=true

[program:ssmtp]
command = /install.sh
